# systemd/niner-miners-api.service
[Unit]
After=syslog.target network.target niner-miners-server.service

[Service]
User=niner-miners
Group=niner-miners
ExecStart=/home/niner-miners/core/actions/api
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=niner-miners-api

[Install]
WantedBy=multi-user.target
